# backups/crontab
# Run the backup at midnight every day
0 0 * * * /backups/backup.sh >> /var/log/cron.log 2>&1
